# dashboard/com.gautam.genai-roadmap-dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.gautam.genai-roadmap-dashboard</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/gautam/Library/Application Support/GenAIRoadmapDashboard/app.py</string>
    <string>--host</string>
    <string>127.0.0.1</string>
    <string>--port</string>
    <string>8765</string>
    <string>--open-browser</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/gautam/Library/Application Support/GenAIRoadmapDashboard</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/gautam/Library/Application Support/GenAIRoadmapDashboard/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/gautam/Library/Application Support/GenAIRoadmapDashboard/launchd.err.log</string>
</dict>
</plist>
